# dockerbase/supervisor/order-srv.conf
[program:order-srv]
command=/data/goapp/mtbsystem/bin/order-srv
process_name=order-srv
autorestart=true
redirect_stderr=true
stderr_logfile=/data/logs/order-srv.err.log
stdout_logfile=/data/logs/order-srv.out.log
user=root